# blocksupply/deploy/blocksupply.service
[Unit]
Description=BlockSupply Flask App (gunicorn)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/blocksupply
EnvironmentFile=/home/ec2-user/blocksupply/.env
ExecStart=/home/ec2-user/blocksupply/venv/bin/gunicorn -w 3 -b 0.0.0.0:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
